# apps/systemd/user/clean-downloads.timer
[Unit]
Description=Run cleaning Downloads daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
