# apps/systemd/user/clean-downloads.service
[Unit]
Description=Remove files older than 7 days from Downloads

[Service]
ExecStart=/usr/bin/find ~/Downloads -atime +7 -delete
